# kindle/keep_ready_to_suspend.conf
start on started lab126
stop on stopping system

exec /mnt/us/wallpaper/keep_ready_to_suspend.sh >> /var/log/keep_ready_to_suspend.log 2>&1
